# deploy/systemd/praxis.service.d/hardening.conf
# systemd hardening drop-in for praxis. Install at
# /etc/systemd/system/praxis.service.d/hardening.conf and `systemctl daemon-reload`.
# Verify with `systemd-analyze security praxis.service`.
#
# Complements the baseline in praxis.service (NoNewPrivileges, ProtectSystem,
# ProtectHome, PrivateTmp, PrivateDevices, ReadWritePaths, StateDirectory), which
# stays in the base unit so an install without this drop-in is still protected;
# this file adds the rest without repeating it (BL-087 dedupe).

[Service]
# Filesystem, process, and kernel isolation.
ProtectKernelTunables=true
ProtectKernelModules=true
ProtectKernelLogs=true
ProtectControlGroups=true
ProtectClock=true
ProtectHostname=true
ProtectProc=invisible
ProcSubset=pid

# Privilege and capability lockdown.
PrivateUsers=true
CapabilityBoundingSet=
AmbientCapabilities=
RestrictSUIDSGID=true
RestrictNamespaces=true
RestrictRealtime=true
LockPersonality=true
MemoryDenyWriteExecute=true
RemoveIPC=true

# System-call surface.
SystemCallFilter=@system-service
SystemCallFilter=~@privileged @resources @obsolete
SystemCallArchitectures=native

# Network: restrict address families; egress is the actuation surface, so keep this
# tight and complement it with host firewalling. A full IP-level lockdown
# (IPAddressDeny=any plus an IPAddressAllow allowlist plus SocketBindDeny=any) is NOT
# preset (BL-087): praxis must reach the operator's fleet over SSH, so a deny-all
# default would brick actuation. A turnkey opt-in ships beside this file as
# network-lockdown.conf.example (ADR-0034): copy it to network-lockdown.conf, scope
# IPAddressAllow to your fleet, and daemon-reload.
RestrictAddressFamilies=AF_UNIX AF_INET AF_INET6

# Writable state files are created group/world-unreadable.
UMask=0077
