# deploy/systemd/praxis.service.d/hardening.conf
# systemd hardening drop-in for praxis. Install at
# /etc/systemd/system/praxis.service.d/hardening.conf and `systemctl daemon-reload`.
# Verify with `systemd-analyze security praxis.service`.

[Service]
# Filesystem and device isolation.
ProtectSystem=strict
ProtectHome=true
PrivateTmp=true
PrivateDevices=true
ProtectKernelTunables=true
ProtectKernelModules=true
ProtectKernelLogs=true
ProtectControlGroups=true
ProtectClock=true
ProtectHostname=true
ProtectProc=invisible

# Privilege and capability lockdown.
NoNewPrivileges=true
CapabilityBoundingSet=
AmbientCapabilities=
RestrictSUIDSGID=true
RestrictNamespaces=true
RestrictRealtime=true
LockPersonality=true
MemoryDenyWriteExecute=true

# System-call surface.
SystemCallFilter=@system-service
SystemCallFilter=~@privileged @resources @obsolete
SystemCallArchitectures=native

# Network: restrict address families; egress is the actuation surface, so keep
# this tight and complement it with host firewalling.
RestrictAddressFamilies=AF_UNIX AF_INET AF_INET6

# Writable state only under the state directory.
StateDirectory=praxis
ReadWritePaths=/var/lib/praxis
UMask=0077
